# com.martijngastkemper.qmd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.martijngastkemper.qmd</string>
	<key>ProgramArguments</key>
	<array>
  	<string>/bin/sh</string>
		<string>-c</string>
		<string>export PATH="$HOME/.nvm/versions/node/$(ls ~/.nvm/versions/node | sort -V | tail -1)/bin:$PATH" &amp;&amp; ~/.bun/bin/qmd mcp --http</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>StandardOutPath</key>
	<string>/tmp/qmd.stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/tmp/qmd.stderr.log</string>
</dict>
</plist>
